# crontab
 */1 * * * /home/theotim/backup.sh >> backup.log
